# ros_catkin_ws/devel/etc/catkin/profile.d/15.rosbash.bash
# generated from rosbash/env-hooks/15.rosbash.bash.em

. "/home/odroid/ros_catkin_ws/src/ros/rosbash/rosbash"
